# deploy/bbc-radio-search.service
[Unit]
Description=BBC Radio Player Podcast Search Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=shaivure
WorkingDirectory=/home/shaivure/bbc-radio/api
Environment=BBC_RADIO_DATA_DIR=/home/shaivure/bbc-radio/data
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/shaivure/bbc-radio/venv/bin/python search_server.py
Restart=always
RestartSec=5
StandardOutput=append:/home/shaivure/bbc-radio/logs/search-server.log
StandardError=append:/home/shaivure/bbc-radio/logs/search-server.log

# Resource limits
MemoryMax=2G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
